# build/crontab
* * * * * cd /var/www && php8 artisan schedule:run >> /dev/stdout 2>&1
