# build/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm7 -F
user=app
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/var/log/php7/stdout.log
stderr_logfile=/var/log/php7/stderr.log
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
user=app
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log
autorestart=false
startretries=0

[program:websockets]
command=php /app/artisan websockets:serve
user=app
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/var/log/websockets/stdout.log
stderr_logfile=/var/log/websockets/stderr.log
autorestart=false
startretries=0

[program:cron]
command=crond -f -L /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
autorestart=true
startretries=5
